# docker/s6/gogs/run
#!/bin/sh
echo 'gogs run'
if test -f ./setup; then
    source ./setup
fi

export USER=git
exec gosu $USER /app/gogs/gogs web
